# supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log

[program:flask]
command=/usr/bin/python3 /home/claude/webapps/flask/app.py
directory=/home/claude/webapps/flask
user=claude
autostart=false
autorestart=true
stderr_logfile=/var/log/flask-error.log
stdout_logfile=/var/log/flask-access.log
environment=FLASK_APP="app.py",FLASK_ENV="development"

[program:ttyd]
# ttyd with auto-login via sshpass - opens directly to shell
command=/usr/bin/ttyd -p 7682 -W -t fontSize=12 sshpass -p claude ssh -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null claude@localhost
autostart=true
autorestart=true
stderr_logfile=/var/log/ttyd-error.log
stdout_logfile=/var/log/ttyd-access.log

